# pc-power-scheduler/usr/share/pc-power-scheduler/pc-power-scheduler.path
[Path]
PathModified=/etc/default/pc-power-scheduler

[Install]
WantedBy=multi-user.target