# pc-power-scheduler/usr/share/pc-power-scheduler/shutdown.timer.tmpl
[Unit]
Description=Shuts the computer down at $SHUTDOWNTIME

[Timer]
Unit=shutdown.service
OnCalendar=$DAYS $SHUTDOWNTIME
AccuracySec=1s
Persistent=true

[Install]
WantedBy=timers.target