# pc-power-scheduler/usr/share/pc-power-scheduler/wol.service.tmpl
[Unit]
Description=Enable Wake On Lan

[Service]
Type=oneshot
ExecStart=ethtool --change ${NICNAME} wol g

[Install]
WantedBy=basic.target